# example/ipshow/ipshow.service
[Unit]
Description=IP Show
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/ipshow

[Install]
WantedBy=multi-user.target

